# supervisor.conf
[program:jukebox]
command  = python -u /home/pi/jukebox/jukebox.py /home/pi/jukebox/cards.yaml
logfile  = /home/pi/jukebox/log.txt
loglevel = info
user     = pi
